# etc/cron.d/lwes-journaller-rotate
* * * * * root /etc/init.d/lwes-journaller rotate >/dev/null 2>&1
