# etc/supervisord.d/lwes-journaller.ini
[program:lwes-journaller]
command=/etc/rc.d/init.d/lwes-journaller foreground
directory=/
autostart=true                 ; start at supervisord start (default: true)
autorestart=true               ; restart at unexpected quit (default: true)
priority=900                  ; the relative start priority (default 999)
environment=LANG="C",TZ="UTC"
